# deployment/slayflashcards-backend-staging.service
[Unit]
Description=SlayFlashcards Backend API (Staging)
After=network.target

[Service]
Type=simple
User=github
Group=github
WorkingDirectory=/opt/slayflashcards/staging/backend
Environment="PATH=/opt/slayflashcards/staging/backend/.venv/bin"
Environment="DATABASE_URL=sqlite:////opt/slayflashcards/staging/backend/slayflashcards-staging.db"
ExecStart=/opt/slayflashcards/staging/backend/.venv/bin/uvicorn api.main_api:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
